# configs/storage-agent.service
[Unit]
Description=Server Storage Management System Agent
Documentation=file:/opt/ssms/agent/README.md
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
EnvironmentFile=-/etc/ssms/storage-agent.env
Environment=SSMS_SERVER_URL=http://127.0.0.1:8080
Environment=SSMS_AGENT_NAME=node01
Environment=SSMS_AGENT_ADDRESS=127.0.0.1
Environment=SSMS_AGENT_DISK=/
Environment=SSMS_AGENT_INTERVAL=30s
ExecStart=/usr/local/bin/storage-agent -server ${SSMS_SERVER_URL} -name ${SSMS_AGENT_NAME} -address ${SSMS_AGENT_ADDRESS} -disk ${SSMS_AGENT_DISK} -interval ${SSMS_AGENT_INTERVAL}
Restart=always
RestartSec=5
NoNewPrivileges=true

[Install]
WantedBy=multi-user.target
